# nanowakeword/rootfs/etc/services.d/nanowakeword/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

MODEL_DIR="$(bashio::config 'model_dir')"
DEFAULT_MODEL="$(bashio::config 'default_model')"
THRESHOLD="$(bashio::config 'threshold')"
TRIGGER_LEVEL="$(bashio::config 'trigger_level')"
REFRACTORY_SECONDS="$(bashio::config 'refractory_seconds')"
VAD_THRESHOLD="$(bashio::config 'vad_threshold')"
CASCADE="$(bashio::config 'cascade')"
GATE_THRESHOLD="$(bashio::config 'gate_threshold')"
HTTP_API="$(bashio::config 'http_api')"
CAPTURE_DETECTIONS="$(bashio::config 'capture_detections')"
DEBUG_LOGGING="$(bashio::config 'debug_logging')"

# The model dir must exist before the server scans it and before the HTTP
# API accepts uploads into it.
mkdir -p "${MODEL_DIR}"

args=(
  "wyoming-nanowakeword"
  "--uri" "tcp://0.0.0.0:10400"
  "--model-dir" "${MODEL_DIR}"
  "--threshold" "${THRESHOLD}"
  "--trigger-level" "${TRIGGER_LEVEL}"
  "--refractory-seconds" "${REFRACTORY_SECONDS}"
  "--vad-threshold" "${VAD_THRESHOLD}"
  "--gate-threshold" "${GATE_THRESHOLD}"
  "--zeroconf" "nanoWakeWord"
)

if [[ -n "${DEFAULT_MODEL}" ]]; then
  args+=("--default-model" "${DEFAULT_MODEL}")
fi

if [[ "${CASCADE}" == "true" ]]; then
  args+=("--cascade")
fi

if [[ "${CAPTURE_DETECTIONS}" == "true" ]]; then
  args+=("--capture-dir" "${MODEL_DIR}/captures")
fi

if [[ "${HTTP_API}" == "true" ]]; then
  args+=("--http-port" "10401")
  if bashio::config.has_value 'http_token'; then
    HTTP_TOKEN="$(bashio::config 'http_token')"
  else
    # Never expose the API unauthenticated: generate a persistent token.
    TOKEN_FILE=/data/http_token
    if [[ ! -s "${TOKEN_FILE}" ]]; then
      python3 -c "import secrets; print(secrets.token_urlsafe(24))" > "${TOKEN_FILE}"
    fi
    HTTP_TOKEN="$(cat "${TOKEN_FILE}")"
    bashio::log.info \
      "HTTP model API token (set the http_token option to override): ${HTTP_TOKEN}"
  fi
  args+=("--http-token" "${HTTP_TOKEN}")
fi

if [[ "${DEBUG_LOGGING}" == "true" ]]; then
  args+=("--debug")
fi

exec "${args[@]}"
